# dist/etc/init.d/flapjack
#!/bin/bash

### BEGIN INIT INFO
# Provides:       flapjack
# Required-Start: $syslog $remote_fs redis-flapjack
# Required-Stop:  $syslog $remote_fs redis-flapjack
# Should-Start:   $local_fs
# Should-Stop:    $local_fs
# Default-Start:  2 3 4 5
# Default-Stop:   0 1 6
# Short-Description:  flapjack - scalable monitoring notification system
# Description:    flapjack - scalable monitoring notification system
### END INIT INFO

# Copyright (c) 2009-2013 Lindsay Holmwood <lindsay@holmwood.id.au>
#
# Boots flapjack (coordinator, processor, notifier, gateways...)

PATH=/opt/flapjack/bin:$PATH

if [ ! $(type flapjack) ]; then
  echo "Error: flapjack isn't in PATH."
  echo "Refusing to do anything!"
  exit 1
fi

# Evaluate command
flapjack server $@

RETVAL=$?
exit $RETVAL

